# deploy/musicweb-backup.timer
[Unit]
Description=Daily MusicWeb backup

[Timer]
OnCalendar=*-*-* 03:15:00
Persistent=true

[Install]
WantedBy=timers.target
