# deploy/musicweb-backup.service
[Unit]
Description=Backup MusicWeb SQLite databases

[Service]
Type=oneshot
User=musicweb
Group=musicweb
ExecStart=/opt/musicweb/app/deploy/backup.sh
